# deploy/alsamixer-web.service
[Unit]
Description=ALSA Mixer Web Interface
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/work/alsamixer-web
ExecStart=/root/work/alsamixer-web/alsamixer-web-wrapper
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
